# user/.config/autostart/startup-sound.desktop
[Desktop Entry]
Name=Startup Sound
Comment=Plays a startup noise.
Exec=play ~/Desktop/startup.mp3
Type=Application
Terminal=false
Hidden=false
